# src/alert-agent.service.in
[Unit]
Description=alert-agent service
After=malamute.service
Requires=malamute.service

[Service]
Type=simple
ExecStart=@prefix@/bin/alert-agent @sysconfdir@/@PACKAGE@/alert-agent.cfg

[Install]
WantedBy=multi-user.target
