# src/bios-agent-alert-generator.service.in
[Unit]
Description=fty-agent-alert-generator service
Requires=network.target malamute.service
After=network.target malamute.service
PartOf=fty.target

[Service]
Type=simple
User=fty
EnvironmentFile=-/usr/share/fty/etc/default/fty
EnvironmentFile=-/usr/share/fty/etc/default/fty__%n.conf
EnvironmentFile=-/etc/default/fty
EnvironmentFile=-/etc/default/fty__%n.conf
ExecStart=@prefix@/bin/fty-agent-alert-generator
Restart=always

[Install]
WantedBy=fty.target
