# src/bios-agent-alert-generator.service.in
[Unit]
Description=bios-agent-alert-generator service
Requires=network.target malamute.service
After=network.target malamute.service
PartOf=bios.target

[Service]
Type=simple
ExecStart=@prefix@/bin/bios-agent-alert-generator

[Install]
WantedBy=bios.target
