# src/bios-agent-alert-generator.service.in
[Unit]
Description=bios-agent-alert-generator service
Requires=network.target malamute.service
After=network.target malamute.service
PartOf=bios.target

[Service]
Type=simple
User=bios
EnvironmentFile=-/usr/share/bios/etc/default/bios
EnvironmentFile=-/usr/share/bios/etc/default/bios__%n.conf
EnvironmentFile=-/etc/default/bios
EnvironmentFile=-/etc/default/bios__%n.conf
ExecStart=@prefix@/bin/bios-agent-alert-generator
Restart=always

[Install]
WantedBy=bios.target
